# bootstrap-v21c.service
[Unit]
Description=Autonomous LLM Agent Substrate (%i)
After=network.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=bootstrap-v21c
WorkingDirectory=/home/bootstrap-v21c/bootstrap
Environment="NODE_EXTRA_CA_CERTS=/usr/local/share/ca-certificates/mitmproxy/mitmproxy.crt"
ExecStart=/home/bootstrap-v21c/bootstrap/run-agent.sh
Restart=always
RestartSec=5s
StandardOutput=inherit
StandardError=inherit
SyslogIdentifier=bootstrap-v21c

[Install]
WantedBy=multi-user.target
